# systemd/kp4pra-morse-id.timer
[Unit]
Description=KP4PRA TNC - blink Morse callsign ID every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=10s

[Install]
WantedBy=timers.target
