# systemd/kp4pra-tnc-agent.service
[Unit]
Description=KP4PRA TNC - Bluetooth NoInputNoOutput Pairing Agent (Just Works)
After=bluetooth.service
Requires=bluetooth.service
PartOf=kp4pra-tnc.target

[Service]
Type=simple
ExecStart=/usr/bin/bt-agent --capability=NoInputNoOutput
Restart=on-failure
# bt-agent ignores SIGTERM when D-Bus is torn down first (reboot race),
# hanging shutdown for the 90s default. It is stateless - kill fast.
TimeoutStopSec=5
RestartSec=5

[Install]
WantedBy=multi-user.target kp4pra-tnc.target
